# crontab
# Run the script every hour
0 * * * * cd /app && /usr/local/bin/python3 polar_test.py >> /var/log/cron.log 2>&1
